# sample-login/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:java]
command=java -jar sample-login.jar

[program:cron]
command=cron -f

[program:nginx]
command=service nginx start

